# files/root/etc/uci-defaults/luci-beardropper
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@bearDropper[-1]
	add ucitrack bearDropper
	set ucitrack.@bearDropper[-1].init=bearDropper
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
